# backup/etc/runit/sv/avahi-daemon/run
#!/bin/sh
# if dbus is enabled wait for it.
if [ -e /run/runit/service/dbus ]; then
	sv check dbus > /dev/null || exit 1
fi
exec avahi-daemon -s
