# backup/etc/NetworkManager/dispatcher.d/0-vpn-down
#!/bin/bash

if [ -z "$(ip a | grep vpn0 | grep DOWN)" ]
then
	echo "VPN is up, everything should be working..." 
else
	echo "VPN is down, reconnecting"
	vpn-reconnect
fi