# backup/etc/runit/sv/clear_cache/run
#!/bin/sh

set -e;
exec 2>&1;

while true; do
  snooze -v -H'*' -M/5 clear_cache;
done;
